# project4-with-file-provision/atbooting.sh
#!/bin/bash
# chkconfig: - 50 50
# description: init file for Nodejs server
sleep 10
cd ~
cd simple-nodejs/
forever start --minUptime 3000 --spinSleepTime 3000  index.js > reboot.logs